# deploy/supervisor.ini
[program:battery_exporter]
command=/usr/bin/python battery_exporter.py
directory=/PATH/TO/prometheus_battery_exporter
autostart=true
autorestart=true
stdout_logfile=/var/log/battery_exporter.stdout
stderr_logfile=/var/log/battery_exporter.stderr
